# playbooks/templates/ansible-config.service.j2
[Unit]
Description=Ansible pull test

[Service]
Type=oneshot
ExecStart={{ shell_command }}